# the_little_init_system_that_could/resources/tollur.service
[Unit]
Description=Nimble and scriptable SMTP proxy

[Service]
Type=simple
ExecStart=/usr/bin/tollur /etc/tollur/conf.ini

[Install]
WantedBy=multi-user.target
